# MedAgent-dev-jan/systemd/medagent-backend.service
[Unit]
Description=MedAgent Backend Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/medagent/sample-deploy/MedAgent-dev-jan
Environment="PATH=/opt/medagent/sample-deploy/MedAgent-dev-jan/.venv/bin"
ExecStart=/opt/medagent/sample-deploy/MedAgent-dev-jan/.venv/bin/uvicorn backend.app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
